# systemd/hillview.timer
[Unit]
Description=hillview timer

[Timer]
AccuracySec=50s
OnBootSec=1000s
OnUnitActiveSec=500s
Unit=hillview.service

[Install]
WantedBy=timers.target


